# .local/src/etc/udev/rules.d/60-mouse.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", RUN+="imwheel -k"
